# resources/etc/sv/zramswap/run
#!/bin/sh
exec /usr/bin/zramctrl start 2>&1
